# Reaktor/etc/systemd/system/Reaktor.service
[Unit]
Description=Reaktor for user %i
After=network.target nss-lookup.target

[Service]
Type=forking
User=%i
ExecStart=/krebs/painload/Reaktor/index
Restart=always

[Install]
WantedBy=multi-user.target
