# Reaktor/etc/systemd/system/Reaktor@.service
[Unit]
Description=Reaktor for user %i
After=network.target nss-lookup.target

[Service]
Type=normal
#TODO - make reaktor path variable
User=%i
ExecStart=/krebs/painload/Reaktor/index
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
