# god/Monitoring/collectd/supervisor.d/collectd.conf
[program:collectd]
command=collectd -C /etc/collectd.conf -f
autorestart=true
