# minikrebs/profiles/krebs/customfiles/etc/rc.local
/usr/bin/genhostid
[ ! -e /etc/tinc/retiolum ] && /usr/bin/autoinfest
exit 0
