# Reaktor/startup/supervisor/Reaktor.conf
[program:Reaktor]
command=/usr/bin/python2.6 IRC/asybot.py
environment=host='irc.freenode.net',target='#krebsco' 
redirect_stderr=true
user=reaktor
directory=/krebs/Reaktor
