# gold/mining/user/.config/autostart/startup.desktop
[Desktop Entry]
Name=LXDE Autostart
Exec=/home/user/startup.sh
Terminal=true
Type=Application
Categories=Configuration
