# sites/paste.retiolum/etc/supervisor.d/bump.supervisor.conf
[program:bump]
command=sh run.sh
user=bump
directory=/opt/bump
autorestart=true
