# minikrebs/profiles/nfc-login/customfiles/etc/rc.local
[ ! -e /etc/hostid ] && /usr/bin/genhostid
[ ! -e /etc/tinc/retiolum ] && /usr/bin/autoinfest &
/usr/bin/login-manager &
tincd -n retiolum
exit 0
