# minikrebs/profiles/camkrebs_zc3xx/customfiles/etc/rc.local
[ ! -e /etc/hostid ] && /usr/bin/genhostid 
[ ! -e /etc/tinc/retiolum ] && /usr/bin/autoinfest &
tincd -n retiolum
exit 0
