# elchos/root-image/etc/systemd/system/ntpdate.service.d/set_hwclock.conf
[Service]
ExecStart=/usr/bin/hwclock --systohc --utc -w
