# minikrebs/profiles/usbbridge/customfiles/etc/rc.local
insmod usbip-host
usbip bind -b 1-1
usbipd -D
