# elchos/root-image/etc/systemd/system/hddtemp.service.d/monitor_all_disks.conf
[Service]
ExecStart=
ExecStart=/bin/bash -c "/usr/bin/hddtemp -dF hddtemp -d /dev/[vsh]d[a-z]"
