# elchos/root-image/etc/systemd/system/collectd.service.d/wait.conf
[Service]
ExecStartPre=/usr/bin/sleep 20
