# god/twinter/supervisor.d/twinter.conf
[program:twinter]
command=/bin/sh /krebs/god/twinter/init.sh
directory=/krebs/god/twinter/
user=nobody
autorestart=true
