# retiolum/scripts/github_listener/github_listener.conf
[program:github_listener]
command=/krebs/retiolum/scripts/github_listener/listen 5432
user=nobody
